# .crontab
SHELL=/bin/bash
DISPLAY=":0"
XAUTHORITY="/home/gregdan3/.Xauthority"
XDG_RUNTIME_DIR="/run/user/1000"
DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
0 12 * * * gopass sync
0 */2 * * * dunstify --appname="pacman" "$(checkupdates --download | wc -l) updates available"
*/5 * * * * curl -s "wttr.in/?format=3" > ~/.weather
